# app/supervisord/base.conf
[supervisord]
nodaemon = true

[program:sshd]
priority = 1
command = /bin/bash -c "/usr/sbin/sshd -D"
