# app/supervisord/grafana.conf

[program:grafana]
priority = 3
command = /bin/bash -c "/etc/init.d/grafana-server start"
