# app/supervisord/influxdb.conf

[program:influxdb]
priority = 2
command = /bin/bash -c "/etc/init.d/influxdb start && sleep 5"
